# clients/pades-token-agent-cross-platform/packaging/linux/maiocchi-pades-agent.service
[Unit]
Description=Maiocchi PAdES Token Agent
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/libexec/maiocchi/maiocchi-pades-token-agent
Restart=on-failure
RestartSec=3
Environment=MAIOCCHI_PORTAL_ORIGIN=https://assinatura.maiocchi.adv.br
Environment=MAIOCCHI_ALLOWED_ORIGINS=https://assinatura.maiocchi.adv.br
EnvironmentFile=-%h/.config/maiocchi/pades-agent.env
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_INET AF_UNIX
IPAddressDeny=any
IPAddressAllow=localhost
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=default.target
